# debian/additional_files/lib/systemd/system/user@.service.d/benchexec.conf
[Service]
# Delegate all controllers such that BenchExec can make use of them
Delegate=yes
